# src/generate_container_packages/templates/debian/postinst.j2
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create configuration directory if it doesn't exist
        mkdir -p "{{ paths.etc }}"

        # Always update env.defaults from template (contains system defaults)
        cp "{{ paths.lib }}/env.template" "{{ paths.etc }}/env.defaults"

        # Create env for user overrides if it doesn't exist (with commented defaults)
        if [ ! -f "{{ paths.etc }}/env" ]; then
            cp "{{ paths.lib }}/env.user-template" "{{ paths.etc }}/env"
        fi

        # Source environment to resolve CONTAINER_DATA_ROOT
        . "{{ paths.etc }}/env.defaults"

        # Create data directories with proper ownership for non-root containers
        # Only set ownership on newly created directories to avoid overwriting user changes
        {% for vol in service.volume_directories %}
        {% if vol.path.startswith('${CONTAINER_DATA_ROOT}') or vol.path.startswith('$CONTAINER_DATA_ROOT') %}
        # Resolve CONTAINER_DATA_ROOT in path
        VOL_PATH=$(echo "{{ vol.path }}" | envsubst)
        {% else %}
        VOL_PATH="{{ vol.path }}"
        {% endif %}
        if [ ! -d "$VOL_PATH" ]; then
            mkdir -p "$VOL_PATH"
            {% if vol.uid is not none and vol.gid is not none %}
            chown {{ vol.uid }}:{{ vol.gid }} "$VOL_PATH"
            {% elif vol.uid is not none %}
            chown {{ vol.uid }} "$VOL_PATH"
            {% endif %}
        fi
        {% endfor %}

        # Only interact with systemd if it's running (not in containers or minimal envs)
        if [ -d /run/systemd/system ]; then
            # Reload systemd to recognize new service
            systemctl daemon-reload

            # Enable service (but don't start automatically - let user configure first)
            # Users can start with: systemctl start {{ package.name }}.service
            systemctl enable {{ package.name }}.service || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
