# src/generate_container_packages/templates/debian/postinst.j2
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create configuration directory if it doesn't exist
        mkdir -p "{{ paths.etc }}"

        # Always update env.defaults from template (contains system defaults)
        cp "{{ paths.lib }}/env.template" "{{ paths.etc }}/env.defaults"

        # Create env for user overrides if it doesn't exist (with commented defaults)
        if [ ! -f "{{ paths.etc }}/env" ]; then
            cp "{{ paths.lib }}/env.user-template" "{{ paths.etc }}/env"
        fi

        # Source environment to resolve CONTAINER_DATA_ROOT
        . "{{ paths.etc }}/env.defaults"

        # Create data directories with proper ownership for non-root containers
        # Only set ownership on newly created directories to avoid overwriting user changes
        {% for vol in service.volume_directories %}
        # Resolve any environment variables in path using shell-native expansion
        eval VOL_PATH=\"{{ vol.path }}\"
        if [ ! -d "$VOL_PATH" ]; then
            mkdir -p "$VOL_PATH"
            {% if vol.uid is not none and vol.gid is not none %}
            chown {{ vol.uid }}:{{ vol.gid }} "$VOL_PATH"
            {% elif vol.uid is not none %}
            chown {{ vol.uid }} "$VOL_PATH"
            {% endif %}
        fi
        {% endfor %}

{% if is_oidc_app %}
        # Generate OIDC client secret if not exists
        OIDC_SECRET_FILE="${CONTAINER_DATA_ROOT}/oidc-secret"
        if [ ! -f "${OIDC_SECRET_FILE}" ]; then
            openssl rand -hex 32 > "${OIDC_SECRET_FILE}"
            chmod 600 "${OIDC_SECRET_FILE}"
        fi
{% endif %}

{% if has_default_data %}
        # Copy default data files to data volume on first install (if destination doesn't exist)
        # Note: Uses first volume's uid/gid for ownership. If no volumes defined, files are root-owned.
        DEFAULT_DATA_SRC="{{ paths.lib }}/default-data"
        if [ -d "$DEFAULT_DATA_SRC" ]; then
            # Find all files using -print0 for robustness with special characters in filenames
            find "$DEFAULT_DATA_SRC" -type f -print0 | while IFS= read -r -d '' src_file; do
                # Get relative path from default-data directory
                rel_path="${src_file#$DEFAULT_DATA_SRC/}"
                # Expand CONTAINER_DATA_ROOT in destination path
                dst_file="${CONTAINER_DATA_ROOT}/${rel_path}"
                dst_dir="$(dirname "$dst_file")"

                # Only copy if destination file doesn't exist (preserve user modifications)
                if [ ! -f "$dst_file" ]; then
                    # Create destination directory if needed
                    mkdir -p "$dst_dir"
                    # Copy file preserving permissions
                    cp -p "$src_file" "$dst_file"
{% if service.volume_directories %}
{% set first_vol = service.volume_directories[0] %}
{% if first_vol.uid is not none and first_vol.gid is not none %}
                    # Set ownership to match container user (file and all parent dirs up to CONTAINER_DATA_ROOT)
                    chown {{ first_vol.uid }}:{{ first_vol.gid }} "$dst_file"
                    # Chown all directories from dst_dir up to CONTAINER_DATA_ROOT
                    current_dir="$dst_dir"
                    while [ "$current_dir" != "$CONTAINER_DATA_ROOT" ] && [ "$current_dir" != "/" ]; do
                        chown {{ first_vol.uid }}:{{ first_vol.gid }} "$current_dir"
                        current_dir="$(dirname "$current_dir")"
                    done
{% endif %}
{% endif %}
                fi
            done
        fi
{% endif %}

        # Only interact with systemd if it's running (not in containers or minimal envs)
        if [ -d /run/systemd/system ]; then
            # Reload systemd to recognize new service
            systemctl daemon-reload

            # Enable service (but don't start automatically - let user configure first)
            # Users can start with: systemctl start {{ package.name }}.service
            systemctl enable {{ package.name }}.service || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
